# deploy/com.agentmail.relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- AgentMail relay as a launchd LaunchAgent (macOS peers).
     install: cp com.agentmail.relay.plist ~/Library/LaunchAgents/
              launchctl load ~/Library/LaunchAgents/com.agentmail.relay.plist
     Ensure the same bearer token as other machines (config.json or AGENTMAIL_TOKEN in EnvironmentVariables). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>            <string>com.agentmail.relay</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>"$HOME/.dotnet/tools/agentmail" serve</string>
    </array>
    <key>RunAtLoad</key>        <true/>
    <key>KeepAlive</key>        <true/>
</dict>
</plist>
